# under_attack_ddos/deployment/uad-l4-analyzer.service
[Unit]
Description=Under Attack L4 SYN Flood Analyzer
After=network.target

[Service]
# Needs root for Scapy/eBPF
User=root
Group=uad-core
WorkingDirectory=/opt/uad
ExecStart=/opt/uad/.venv/bin/python3 -u layer4/syn_flood_analyzer.py --config config/thresholds.yaml --ebpf --json

# Hardening
ProtectSystem=full
ProtectHome=yes
PrivateTmp=yes

# Restart Logic
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
